# deploy/systemd/nodebridge.service
[Unit]
Description=NodeBridge proxy control plane
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/nodebridged -config /etc/nodebridge/config.json
WorkingDirectory=/var/lib/nodebridge
Restart=on-failure
RestartSec=3
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target

